# live-build-config/config/includes.chroot/etc/systemd/system/nakumi-test-automation.service
[Unit]
Description=NakumiOS Automated UI Test - Screenshot and System Info
After=graphical.target greetd.service
Wants=graphical.target

[Service]
Type=oneshot
User=nakumi
ExecStartPre=/bin/sleep 15
ExecStart=/usr/bin/nakumi-test-automation
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
